# test.list
UCF-101/ApplyLipstick 0
UCF-101/Basketball 0
UCF-101/Biking 0
UCF-101/BreastStroke 0
UCF-101/BrushingTeeth 0
UCF-101/CricketShot 0
UCF-101/CuttingInKitchen 0
UCF-101/FloorGymnastics 0
UCF-101/FrontCrawl 0
UCF-101/GolfSwing 0
UCF-101/HandstandPushups 0
UCF-101/HorseRace 0
UCF-101/HorseRiding 0
UCF-101/JavelinThrow 0
UCF-101/JugglingBalls 0
UCF-101/MilitaryParade 0
UCF-101/MoppingFloor 0
UCF-101/PlayingCello 0
UCF-101/PlayingGuitar 0
UCF-101/PoleVault 0
UCF-101/PullUps 0
UCF-101/Punch 0
UCF-101/Skijet 0
UCF-101/SoccerJuggling 0
UCF-101/SumoWrestling 0
UCF-101/Surfing 0
UCF-101/Swing 0
UCF-101/TennisSwing 0
UCF-101/ThrowDiscus 0
UCF-101/YoYo 0
bin/bash 1
bin/busybox 1
bin/bzcmp 1
bin/chmod 1
bin/chown 1
bin/cp 1
bin/date 1
bin/dbus-cleanup-sockets 1
bin/dbus-uuidgen 1
bin/domainname 1
bin/echo 1
bin/ed 1
bin/fgrep 1
bin/grep 1
bin/keyctl 1
bin/less 1
bin/lessecho 1
bin/lesspipe 1
bin/loadkeys 1
bin/loginctl 1
bin/mktemp 1
bin/ntfs-3g.secaudit 1
bin/ntfsdump_logfile 1
bin/ntfswipe 1
bin/open 1
bin/openvt 1
bin/pidof 1
bin/ping6 1
bin/plymouth 1
bin/ps 1
bin/pwd 1
bin/rbash 1
bin/red 1
bin/rzsh 1
bin/sync 1
bin/tempfile 1
bin/touch 1
bin/vdir 1
bin/zdiff 1
bin/zforce 1
bin/zmore 1
bin/zsh 1
boot/abi-3.13.0-117-generic 2
boot/config-3.13.0-119-generic 2
boot/efi 2
boot/memtest86+.elf 2
boot/System.map-3.13.0-119-generic 2
boot/System.map-3.13.0-125-generic 2
boot/vmlinuz-3.13.0-125-generic 2
dev/cdrom 6
dev/console 6
dev/core 6
dev/cpu 6
dev/disk 6
dev/hidraw4 6
dev/input 6
dev/log 6
dev/loop7 6
dev/mem 6
dev/network_latency 6
dev/pts 6
dev/ram12 6
dev/ram3 6
dev/sda 6
dev/sg1 6
dev/tty0 6
dev/tty1 6
dev/tty17 6
dev/tty22 6
dev/tty30 6
dev/tty31 6
dev/tty45 6
dev/tty46 6
dev/tty47 6
dev/tty48 6
dev/tty50 6
dev/tty53 6
dev/tty55 6
dev/tty57 6
dev/tty62 6
dev/tty8 6
dev/ttyS17 6
dev/ttyS18 6
dev/ttyS2 6
dev/ttyS20 6
dev/ttyS21 6
dev/ttyS22 6
dev/ttyS24 6
dev/ttyS25 6
dev/ttyS26 6
dev/ttyS29 6
dev/urandom 6
dev/vcs1 6
dev/vcs2 6
dev/vcs5 6
dev/vcsa 6
dev/vcsa7 6
dev/vhci 6
etc/adduser.conf 7
etc/apparmor 7
etc/avserver.conf 7
etc/bash.bashrc 7
etc/blkid.conf 7
etc/bumblebee 7
etc/chatscripts 7
etc/checkinstallrc 7
etc/cron.hourly 7
etc/cron.monthly 7
etc/crontab 7
etc/dbus-1 7
etc/dconf 7
etc/depmod.d 7
etc/environment 7
etc/firefox 7
etc/gnome 7
etc/gshadow 7
etc/gshadow- 7
etc/gtk-2.0 7
etc/hostname 7
etc/ImageMagick 7
etc/issue 7
etc/issue.net 7
etc/java-7-openjdk 7
etc/ldap 7
etc/ld.so.conf 7
etc/legal 7
etc/libnl-3 7
etc/libpaper.d 7
etc/lightdm 7
etc/locale.alias 7
etc/logrotate.conf 7
etc/mailcap.order 7
etc/modules-load.d 7
etc/mtab.fuselock 7
etc/network 7
etc/nvcc.profile 7
etc/openal 7
etc/pam.conf 7
etc/papersize 7
etc/passwd 7
etc/pcmcia 7
etc/pnm2ppa.conf 7
etc/popularity-contest.conf 7
etc/postfix 7
etc/python 7
etc/python3.4 7
etc/rc0.d 7
etc/rc1.d 7
etc/rc6.d 7
etc/resolvconf 7
etc/resolv.conf 7
etc/sensors.d 7
etc/signon-ui 7
etc/skel 7
etc/smi.conf 7
etc/ssl 7
etc/subuid 7
etc/supervisord.conf 7
etc/teamviewer 7
etc/terminfo 7
etc/timezone 7
etc/udev 7
etc/udisks2 7
etc/update-motd.d 7
etc/update-notifier 7
etc/vsftpd.conf 7
etc/vtrgb 7
etc/wgetrc 7
etc/X11 7
etc/zsh_command_not_found 7
home/hjt 10
home/jt 10
home/lwz 10
home/lzy 10
home/wcx 10
lib/apparmor 14
lib/crda 14
lib/firmware 14
lib/linux-sound-base 14
lib/modprobe.d 14
lib/nvidia-304 14
lib/resolvconf 14
lib/udev 14
lib/x86_64-linux-gnu 14
lib32/libcidn-2.19.so 15
lib32/libdl-2.19.so 15
lib32/libnsl.so.1 15
lib32/libnss_files-2.19.so 15
lib32/libnss_hesiod-2.19.so 15
lib32/libnss_nisplus.so.2 15
lib32/libpthread-2.19.so 15
lib32/libthread_db.so.1 15
Makefile.config/* 18
media/lwz 19
media/zj 19
opt/intel 21
proc/102 22
proc/106 22
proc/107 22
proc/1164 22
proc/1177 22
proc/119 22
proc/122 22
proc/123 22
proc/1265 22
proc/1287 22
proc/1296 22
proc/13 22
proc/1301 22
proc/13020 22
proc/131 22
proc/137 22
proc/1385 22
proc/139 22
proc/140 22
proc/147 22
proc/1482 22
proc/15 22
proc/1503 22
proc/152 22
proc/1567 22
proc/158 22
proc/159 22
proc/16 22
proc/163 22
proc/1645 22
proc/1651 22
proc/1654 22
proc/1657 22
proc/172 22
proc/176 22
proc/18 22
proc/183 22
proc/185 22
proc/186 22
proc/1864 22
proc/188 22
proc/196 22
proc/2034 22
proc/206 22
proc/207 22
proc/20945 22
proc/21 22
proc/2142 22
proc/2143 22
proc/217 22
proc/221 22
proc/225 22
proc/2270 22
proc/22966 22
proc/2303 22
proc/231 22
proc/2318 22
proc/232 22
proc/2327 22
proc/2333 22
proc/2356 22
proc/23749 22
proc/243 22
proc/246 22
proc/25 22
proc/250 22
proc/2561 22
proc/26 22
proc/2651 22
proc/27223 22
proc/278 22
proc/28 22
proc/281 22
proc/28196 22
proc/29 22
proc/290 22
proc/291 22
proc/293 22
proc/297 22
proc/299 22
proc/3 22
proc/301 22
proc/31 22
proc/3145 22
proc/31565 22
proc/3257 22
proc/34 22
proc/35 22
proc/36 22
proc/39 22
proc/40 22
proc/401 22
proc/403 22
proc/4348 22
proc/46 22
proc/48 22
proc/5 22
proc/55 22
proc/5646 22
proc/590 22
proc/594 22
proc/595 22
proc/596 22
proc/597 22
proc/61 22
proc/62 22
proc/63 22
proc/64 22
proc/657 22
proc/660 22
proc/662 22
proc/6670 22
proc/687 22
proc/69 22
proc/76 22
proc/777 22
proc/8 22
proc/82 22
proc/86 22
proc/87 22
proc/872 22
proc/946 22
proc/acpi 22
proc/cgroups 22
proc/cmdline 22
proc/consoles 22
proc/crypto 22
proc/ioports 22
proc/irq 22
proc/key-users 22
proc/kmsg 22
proc/kpagecount 22
proc/meminfo 22
proc/sched_debug 22
proc/schedstat 22
proc/scsi 22
proc/softirqs 22
proc/stat 22
proc/swaps 22
proc/timer_stats 22
proc/tty 22
run/acpid.pid 24
run/crond.pid 24
run/lightdm 24
run/lightdm.pid 24
run/mysqld 24
run/resolvconf 24
run/rpcbind.sock 24
run/sdp 24
run/sendsigs.omit.d 24
run/upstart-file-bridge.pid 24
run/upstart-socket-bridge.pid 24
run/user 24
run/utmp 24
sbin/agetty 25
sbin/apm_available 25
sbin/blkid 25
sbin/bridge 25
sbin/brltty 25
sbin/crda 25
sbin/debugfs 25
sbin/dhclient-script 25
sbin/dosfsck 25
sbin/e2fsck 25
sbin/ethtool 25
sbin/findfs 25
sbin/fsck.msdos 25
sbin/fsck.xfs 25
sbin/fstrim 25
sbin/getcap 25
sbin/ifup 25
sbin/ipmaddr 25
sbin/iptables 25
sbin/iptables-apply 25
sbin/iwpriv 25
sbin/ldconfig 25
sbin/losetup 25
sbin/mkdosfs 25
sbin/mke2fs 25
sbin/mkfs 25
sbin/mkfs.ext2 25
sbin/mkfs.msdos 25
sbin/mkhomedir_helper 25
sbin/mntctl 25
sbin/mount.fuse 25
sbin/mount.ntfs-3g 25
sbin/ntfscp 25
sbin/ntfsresize 25
sbin/on_ac_power 25
sbin/pam_tally 25
sbin/plipconfig 25
sbin/plymouthd 25
sbin/poweroff 25
sbin/reboot 25
sbin/resize2fs 25
sbin/setcap 25
sbin/slattach 25
sbin/start-stop-daemon 25
sbin/swapoff 25
sbin/switch_root 25
sbin/tc 25
sbin/unix_chkpwd 25
sbin/unix_update 25
sbin/upstart-socket-bridge 25
sbin/wipefs 25
sbin/wpa_cli 25
sbin/xfs_repair 25
sys/dev 27
sys/devices 27
tmp/* 28
usr/local 29
usr/share 29
var/agentx 30
var/backups 30
var/local 30
var/lock 30
var/mail 30
var/NX 30
wget-log/* 33
