# train.list
UCF-101/ApplyEyeMakeup 0
UCF-101/Archery 0
UCF-101/BabyCrawling 0
UCF-101/BalanceBeam 0
UCF-101/BandMarching 0
UCF-101/BaseballPitch 0
UCF-101/BasketballDunk 0
UCF-101/BenchPress 0
UCF-101/Billiards 0
UCF-101/BlowDryHair 0
UCF-101/BlowingCandles 0
UCF-101/BodyWeightSquats 0
UCF-101/Bowling 0
UCF-101/BoxingPunchingBag 0
UCF-101/BoxingSpeedBag 0
UCF-101/CleanAndJerk 0
UCF-101/CliffDiving 0
UCF-101/CricketBowling 0
UCF-101/Diving 0
UCF-101/Drumming 0
UCF-101/Fencing 0
UCF-101/FieldHockeyPenalty 0
UCF-101/FrisbeeCatch 0
UCF-101/Haircut 0
UCF-101/Hammering 0
UCF-101/HammerThrow 0
UCF-101/HandstandWalking 0
UCF-101/HeadMassage 0
UCF-101/HighJump 0
UCF-101/HulaHoop 0
UCF-101/IceDancing 0
UCF-101/JumpingJack 0
UCF-101/JumpRope 0
UCF-101/Kayaking 0
UCF-101/Knitting 0
UCF-101/LongJump 0
UCF-101/Lunges 0
UCF-101/Mixing 0
UCF-101/Nunchucks 0
UCF-101/ParallelBars 0
UCF-101/PizzaTossing 0
UCF-101/PlayingDaf 0
UCF-101/PlayingDhol 0
UCF-101/PlayingFlute 0
UCF-101/PlayingPiano 0
UCF-101/PlayingSitar 0
UCF-101/PlayingTabla 0
UCF-101/PlayingViolin 0
UCF-101/PommelHorse 0
UCF-101/PushUps 0
UCF-101/Rafting 0
UCF-101/RockClimbingIndoor 0
UCF-101/RopeClimbing 0
UCF-101/Rowing 0
UCF-101/SalsaSpin 0
UCF-101/ShavingBeard 0
UCF-101/Shotput 0
UCF-101/SkateBoarding 0
UCF-101/Skiing 0
UCF-101/SkyDiving 0
UCF-101/SoccerPenalty 0
UCF-101/StillRings 0
UCF-101/TableTennisShot 0
UCF-101/TaiChi 0
UCF-101/TrampolineJumping 0
UCF-101/Typing 0
UCF-101/UnevenBars 0
UCF-101/VolleyballSpiking 0
UCF-101/WalkingWithDog 0
UCF-101/WallPushups 0
UCF-101/WritingOnBoard 0
bin/bunzip2 1
bin/bzcat 1
bin/bzdiff 1
bin/bzegrep 1
bin/bzexe 1
bin/bzfgrep 1
bin/bzgrep 1
bin/bzip2 1
bin/bzip2recover 1
bin/bzless 1
bin/bzmore 1
bin/cat 1
bin/chacl 1
bin/chgrp 1
bin/chvt 1
bin/cpio 1
bin/dash 1
bin/dbus-daemon 1
bin/dd 1
bin/df 1
bin/dir 1
bin/dmesg 1
bin/dnsdomainname 1
bin/dumpkeys 1
bin/efibootmgr 1
bin/egrep 1
bin/false 1
bin/fgconsole 1
bin/findmnt 1
bin/fuser 1
bin/fusermount 1
bin/getfacl 1
bin/gunzip 1
bin/gzexe 1
bin/gzip 1
bin/hostname 1
bin/ip 1
bin/kbd_mode 1
bin/kill 1
bin/kmod 1
bin/lessfile 1
bin/lesskey 1
bin/ln 1
bin/login 1
bin/lowntfs-3g 1
bin/ls 1
bin/lsblk 1
bin/lsmod 1
bin/mkdir 1
bin/mknod 1
bin/more 1
bin/mount 1
bin/mountpoint 1
bin/mt 1
bin/mt-gnu 1
bin/mv 1
bin/nano 1
bin/nc 1
bin/nc.openbsd 1
bin/netcat 1
bin/netstat 1
bin/nisdomainname 1
bin/ntfs-3g 1
bin/ntfs-3g.probe 1
bin/ntfs-3g.usermap 1
bin/ntfscat 1
bin/ntfsck 1
bin/ntfscluster 1
bin/ntfscmp 1
bin/ntfsfix 1
bin/ntfsinfo 1
bin/ntfsls 1
bin/ntfsmftalloc 1
bin/ntfsmove 1
bin/ntfstruncate 1
bin/ping 1
bin/plymouth-upstart-bridge 1
bin/readlink 1
bin/rm 1
bin/rmdir 1
bin/rnano 1
bin/running-in-container 1
bin/run-parts 1
bin/sed 1
bin/setfacl 1
bin/setfont 1
bin/setupcon 1
bin/sh 1
bin/sh.distrib 1
bin/sleep 1
bin/ss 1
bin/static-sh 1
bin/stty 1
bin/su 1
bin/tailf 1
bin/tar 1
bin/true 1
bin/udevadm 1
bin/ulockmgr_server 1
bin/umount 1
bin/uname 1
bin/uncompress 1
bin/unicode_start 1
bin/which 1
bin/whiptail 1
bin/ypdomainname 1
bin/zcat 1
bin/zcmp 1
bin/zegrep 1
bin/zfgrep 1
bin/zgrep 1
bin/zless 1
bin/znew 1
bin/zsh5 1
boot/abi-3.13.0-119-generic 2
boot/abi-3.13.0-125-generic 2
boot/config-3.13.0-117-generic 2
boot/config-3.13.0-125-generic 2
boot/grub 2
boot/initrd.img-3.13.0-117-generic 2
boot/initrd.img-3.13.0-119-generic 2
boot/initrd.img-3.13.0-125-generic 2
boot/memtest86+.bin 2
boot/memtest86+_multiboot.bin 2
boot/System.map-3.13.0-117-generic 2
boot/vmlinuz-3.13.0-117-generic 2
boot/vmlinuz-3.13.0-119-generic 2
cdrom/* 3
core/* 4
data/tp 5
dev/autofs 6
dev/block 6
dev/bsg 6
dev/btrfs-control 6
dev/bus 6
dev/char 6
dev/cpu_dma_latency 6
dev/cuse 6
dev/dri 6
dev/ecryptfs 6
dev/fb0 6
dev/fb1 6
dev/fd 6
dev/full 6
dev/fuse 6
dev/hidraw0 6
dev/hidraw1 6
dev/hidraw2 6
dev/hidraw3 6
dev/hpet 6
dev/kmsg 6
dev/kvm 6
dev/loop0 6
dev/loop1 6
dev/loop2 6
dev/loop3 6
dev/loop4 6
dev/loop5 6
dev/loop6 6
dev/loop-control 6
dev/mapper 6
dev/mcelog 6
dev/net 6
dev/network_throughput 6
dev/null 6
dev/port 6
dev/ppp 6
dev/psaux 6
dev/ptmx 6
dev/ptp0 6
dev/ptp1 6
dev/ram0 6
dev/ram1 6
dev/ram10 6
dev/ram11 6
dev/ram13 6
dev/ram14 6
dev/ram15 6
dev/ram2 6
dev/ram4 6
dev/ram5 6
dev/ram6 6
dev/ram7 6
dev/ram8 6
dev/ram9 6
dev/random 6
dev/rfkill 6
dev/rtc 6
dev/rtc0 6
dev/sda1 6
dev/sda2 6
dev/sda3 6
dev/sda4 6
dev/sg0 6
dev/shm 6
dev/snapshot 6
dev/snd 6
dev/sr0 6
dev/stderr 6
dev/stdin 6
dev/stdout 6
dev/tty 6
dev/tty10 6
dev/tty11 6
dev/tty12 6
dev/tty13 6
dev/tty14 6
dev/tty15 6
dev/tty16 6
dev/tty18 6
dev/tty19 6
dev/tty2 6
dev/tty20 6
dev/tty21 6
dev/tty23 6
dev/tty24 6
dev/tty25 6
dev/tty26 6
dev/tty27 6
dev/tty28 6
dev/tty29 6
dev/tty3 6
dev/tty32 6
dev/tty33 6
dev/tty34 6
dev/tty35 6
dev/tty36 6
dev/tty37 6
dev/tty38 6
dev/tty39 6
dev/tty4 6
dev/tty40 6
dev/tty41 6
dev/tty42 6
dev/tty43 6
dev/tty44 6
dev/tty49 6
dev/tty5 6
dev/tty51 6
dev/tty52 6
dev/tty54 6
dev/tty56 6
dev/tty58 6
dev/tty59 6
dev/tty6 6
dev/tty60 6
dev/tty61 6
dev/tty63 6
dev/tty7 6
dev/tty9 6
dev/ttyprintk 6
dev/ttyS0 6
dev/ttyS1 6
dev/ttyS10 6
dev/ttyS11 6
dev/ttyS12 6
dev/ttyS13 6
dev/ttyS14 6
dev/ttyS15 6
dev/ttyS16 6
dev/ttyS19 6
dev/ttyS23 6
dev/ttyS27 6
dev/ttyS28 6
dev/ttyS3 6
dev/ttyS30 6
dev/ttyS31 6
dev/ttyS4 6
dev/ttyS5 6
dev/ttyS6 6
dev/ttyS7 6
dev/ttyS8 6
dev/ttyS9 6
dev/uhid 6
dev/uinput 6
dev/vcs 6
dev/vcs3 6
dev/vcs4 6
dev/vcs6 6
dev/vcs7 6
dev/vcsa1 6
dev/vcsa2 6
dev/vcsa3 6
dev/vcsa4 6
dev/vcsa5 6
dev/vcsa6 6
dev/vga_arbiter 6
dev/vhost-net 6
dev/zero 6
etc/acpi 7
etc/akonadi 7
etc/aliases 7
etc/aliases.db 7
etc/alternatives 7
etc/anacrontab 7
etc/apg.conf 7
etc/apm 7
etc/apparmor.d 7
etc/apport 7
etc/apt 7
etc/aptdaemon 7
etc/at-spi2 7
etc/avahi 7
etc/axelrc 7
etc/bash_completion 7
etc/bash_completion.d 7
etc/bindresvport.blacklist 7
etc/blkid.tab 7
etc/bluetooth 7
etc/bonobo-activation 7
etc/brlapi.key 7
etc/brltty 7
etc/brltty.conf 7
etc/ca-certificates 7
etc/ca-certificates.conf 7
etc/calendar 7
etc/click 7
etc/colord.conf 7
etc/compizconfig 7
etc/ConsoleKit 7
etc/console-setup 7
etc/cracklib 7
etc/cron.d 7
etc/cron.daily 7
etc/cron.weekly 7
etc/crypttab 7
etc/cups 7
etc/cupshelpers 7
etc/debconf.conf 7
etc/debian_version 7
etc/default 7
etc/deluser.conf 7
etc/dhcp 7
etc/dictionaries-common 7
etc/dkms 7
etc/dnsmasq.d 7
etc/doc-base 7
etc/dpkg 7
etc/drirc 7
etc/emacs 7
etc/exports 7
etc/exports~ 7
etc/fonts 7
etc/fstab 7
etc/fstab.d 7
etc/ftpusers 7
etc/fuse.conf 7
etc/gai.conf 7
etc/gconf 7
etc/gdb 7
etc/gdm 7
etc/ghostscript 7
etc/gimp 7
etc/gnashpluginrc 7
etc/gnashrc 7
etc/gnashthumbnailrc 7
etc/gnome-app-install 7
etc/gnome-settings-daemon 7
etc/gnome-system-tools 7
etc/gnome-vfs-2.0 7
etc/groff 7
etc/group 7
etc/group- 7
etc/grub.d 7
etc/gssapi_mech.conf 7
etc/gtk-3.0 7
etc/gtkmathview 7
etc/gunicorn.d 7
etc/hdparm.conf 7
etc/host.conf 7
etc/hosts 7
etc/hosts.allow 7
etc/hosts.deny 7
etc/hotplug 7
etc/hp 7
etc/idmapd.conf 7
etc/ifplugd 7
etc/init 7
etc/init.d 7
etc/initramfs-tools 7
etc/inputrc 7
etc/insserv 7
etc/insserv.conf 7
etc/insserv.conf.d 7
etc/iproute2 7
etc/kbd 7
etc/kde4rc 7
etc/kernel 7
etc/kernel-img.conf 7
etc/kerneloops.conf 7
etc/ksysguarddrc 7
etc/kubuntu-default-settings 7
etc/ld.so.cache 7
etc/ld.so.conf.bak 7
etc/ld.so.conf.d 7
etc/libao.conf 7
etc/libaudit.conf 7
etc/libreoffice 7
etc/lintianrc 7
etc/localtime 7
etc/logcheck 7
etc/login.defs 7
etc/logrotate.d 7
etc/lsb-release 7
etc/ltrace.conf 7
etc/magic 7
etc/magic.mime 7
etc/mail 7
etc/mailcap 7
etc/manpath.config 7
etc/matplotlibrc 7
etc/mercurial 7
etc/mime.types 7
etc/mke2fs.conf 7
etc/modprobe.d 7
etc/modules 7
etc/mono 7
etc/mtab 7
etc/mtools.conf 7
etc/mysql 7
etc/nanorc 7
etc/netconfig 7
etc/netscsid.conf 7
etc/NetworkManager 7
etc/networks 7
etc/newt 7
etc/nsswitch.conf 7
etc/NX 7
etc/obex-data-server 7
etc/odbc.ini 7
etc/odbcinst.ini 7
etc/ofono 7
etc/OpenCL 7
etc/openmpi 7
etc/opt 7
etc/os-release 7
etc/pam.d 7
etc/passwd- 7
etc/perl 7
etc/pki 7
etc/pm 7
etc/polkit-1 7
etc/ppp 7
etc/printcap 7
etc/profile 7
etc/profile~ 7
etc/profile.d 7
etc/protocols 7
etc/pulse 7
etc/purple 7
etc/python2.7 7
etc/python3 7
etc/rarfiles.lst 7
etc/rc2.d 7
etc/rc3.d 7
etc/rc4.d 7
etc/rc5.d 7
etc/rc.local 7
etc/rcS.d 7
etc/request-key.conf 7
etc/request-key.d 7
etc/rmt 7
etc/rpc 7
etc/rsyslog.conf 7
etc/rsyslog.d 7
etc/rygel.conf 7
etc/samba 7
etc/sane.d 7
etc/screenrc 7
etc/securetty 7
etc/security 7
etc/selinux 7
etc/sensors3.conf 7
etc/services 7
etc/sgml 7
etc/shadow 7
etc/shadow- 7
etc/shadowsocks.json 7
etc/shells 7
etc/signond.conf 7
etc/snmp 7
etc/snmp-mibs-downloader 7
etc/sound 7
etc/spamassassin 7
etc/speech-dispatcher 7
etc/ssh 7
etc/subgid 7
etc/subgid- 7
etc/subuid- 7
etc/sudoers 7
etc/sudoers.d 7
etc/supervisord.conf~ 7
etc/supervisord.d 7
etc/sysctl.conf 7
etc/sysctl.d 7
etc/systemd 7
etc/system-image 7
etc/texmf 7
etc/thermald 7
etc/thunderbird 7
etc/timidity 7
etc/ts.conf 7
etc/ucf.conf 7
etc/ufw 7
etc/updatedb.conf 7
etc/update-manager 7
etc/UPower 7
etc/upstart-xsessions 7
etc/usb_modeswitch.conf 7
etc/vdpau_wrapper.cfg 7
etc/vim 7
etc/vnc 7
etc/vulkan 7
etc/wildmidi 7
etc/wireshark 7
etc/wodim.conf 7
etc/wpa_supplicant 7
etc/xdg 7
etc/xfce4 7
etc/xml 7
etc/xrdp 7
etc/zsh 7
faster_rcnn_models.tgz/* 8
fast_rcnn_models.tgz/* 9
home/amax 10
home/cdn 10
home/chc 10
home/dy 10
home/lsh 10
home/qlq 10
home/xrj 10
home/zj 10
home/zjq 10
initrd.img/* 11
initrd.img.old/* 12
intel_mkl.conf/* 13
lib/brltty 14
lib/cpp 14
lib/hdparm 14
lib/i386-linux-gnu 14
lib/ifupdown 14
lib/init 14
lib/klibc-gLiulUM5C1Zpwc25rCxX8UZ6S-s.so 14
lib/ld-linux.so.2 14
lib/libhandle.so.1 14
lib/libhandle.so.1.0.3 14
lib/libip4tc.so.0 14
lib/libip4tc.so.0.1.0 14
lib/libip6tc.so.0 14
lib/libip6tc.so.0.1.0 14
lib/libiptc.so.0 14
lib/libiptc.so.0.0.0 14
lib/libxtables.so.10 14
lib/libxtables.so.10.0.0 14
lib/lsb 14
lib/modules 14
lib/modules-load.d 14
lib/plymouth 14
lib/recovery-mode 14
lib/systemd 14
lib/terminfo 14
lib/ufw 14
lib/xtables 14
lib32/ld-2.19.so 15
lib32/ld-linux.so.2 15
lib32/libanl-2.19.so 15
lib32/libanl.so.1 15
lib32/libBrokenLocale-2.19.so 15
lib32/libBrokenLocale.so.1 15
lib32/libc-2.19.so 15
lib32/libcidn.so.1 15
lib32/libcrypt-2.19.so 15
lib32/libcrypt.so.1 15
lib32/libc.so.6 15
lib32/libdl.so.2 15
lib32/libm-2.19.so 15
lib32/libmemusage.so 15
lib32/libm.so.6 15
lib32/libnsl-2.19.so 15
lib32/libnss_compat-2.19.so 15
lib32/libnss_compat.so.2 15
lib32/libnss_dns-2.19.so 15
lib32/libnss_dns.so.2 15
lib32/libnss_files.so.2 15
lib32/libnss_hesiod.so.2 15
lib32/libnss_nis-2.19.so 15
lib32/libnss_nisplus-2.19.so 15
lib32/libnss_nis.so.2 15
lib32/libpcprofile.so 15
lib32/libpthread.so.0 15
lib32/libresolv-2.19.so 15
lib32/libresolv.so.2 15
lib32/librt-2.19.so 15
lib32/librt.so.1 15
lib32/libSegFault.so 15
lib32/libthread_db-1.0.so 15
lib32/libutil-2.19.so 15
lib32/libutil.so.1 15
lib64/ld-linux-x86-64.so.2 16
lost+found/* 17
media/amax 19
media/cdn 19
media/dy 19
media/lsh 19
media/nomachine 19
media/wcx 19
mnt/dy 20
opt/click.ubuntu.com 21
opt/ffmpeg 21
opt/google 21
opt/OpenBLAS 21
opt/sublime_text 21
opt/teamviewer 21
proc/1 22
proc/10 22
proc/100 22
proc/101 22
proc/1015 22
proc/103 22
proc/1034 22
proc/105 22
proc/108 22
proc/11 22
proc/110 22
proc/111 22
proc/1114 22
proc/1115 22
proc/112 22
proc/1124 22
proc/1129 22
proc/113 22
proc/11312 22
proc/1138 22
proc/115 22
proc/117 22
proc/118 22
proc/12 22
proc/120 22
proc/1202 22
proc/1203 22
proc/121 22
proc/124 22
proc/125 22
proc/126 22
proc/127 22
proc/128 22
proc/1283 22
proc/129 22
proc/1295 22
proc/13018 22
proc/1311 22
proc/132 22
proc/133 22
proc/134 22
proc/136 22
proc/138 22
proc/13818 22
proc/13923 22
proc/1397 22
proc/14 22
proc/141 22
proc/142 22
proc/1422 22
proc/143 22
proc/1435 22
proc/144 22
proc/145 22
proc/146 22
proc/148 22
proc/149 22
proc/1497 22
proc/150 22
proc/151 22
proc/1519 22
proc/153 22
proc/154 22
proc/1546 22
proc/1549 22
proc/155 22
proc/1554 22
proc/156 22
proc/157 22
proc/1571 22
proc/1603 22
proc/161 22
proc/162 22
proc/164 22
proc/1646 22
proc/1647 22
proc/1650 22
proc/1652 22
proc/1653 22
proc/1655 22
proc/1656 22
proc/166 22
proc/1661 22
proc/167 22
proc/168 22
proc/169 22
proc/17 22
proc/171 22
proc/173 22
proc/174 22
proc/1762 22
proc/177 22
proc/1772 22
proc/1774 22
proc/178 22
proc/179 22
proc/181 22
proc/1814 22
proc/182 22
proc/184 22
proc/187 22
proc/189 22
proc/19 22
proc/191 22
proc/192 22
proc/193 22
proc/194 22
proc/195 22
proc/197 22
proc/198 22
proc/199 22
proc/2 22
proc/20 22
proc/200 22
proc/2009 22
proc/201 22
proc/2013 22
proc/2014 22
proc/202 22
proc/203 22
proc/204 22
proc/2045 22
proc/205 22
proc/2057 22
proc/2059 22
proc/208 22
proc/2083 22
proc/209 22
proc/210 22
proc/211 22
proc/212 22
proc/213 22
proc/2133 22
proc/214 22
proc/215 22
proc/216 22
proc/218 22
proc/219 22
proc/22 22
proc/220 22
proc/222 22
proc/2222 22
proc/2226 22
proc/2229 22
proc/223 22
proc/2234 22
proc/2239 22
proc/224 22
proc/226 22
proc/227 22
proc/228 22
proc/229 22
proc/2292 22
proc/23 22
proc/230 22
proc/2320 22
proc/2321 22
proc/2322 22
proc/2323 22
proc/233 22
proc/234 22
proc/235 22
proc/236 22
proc/2363 22
proc/237 22
proc/2373 22
proc/2375 22
proc/238 22
proc/24 22
proc/240 22
proc/241 22
proc/2413 22
proc/244 22
proc/245 22
proc/247 22
proc/249 22
proc/251 22
proc/25129 22
proc/252 22
proc/253 22
proc/254 22
proc/255 22
proc/256 22
proc/257 22
proc/2662 22
proc/2671 22
proc/269 22
proc/27 22
proc/272 22
proc/27277 22
proc/273 22
proc/274 22
proc/275 22
proc/276 22
proc/277 22
proc/279 22
proc/280 22
proc/282 22
proc/283 22
proc/284 22
proc/285 22
proc/286 22
proc/287 22
proc/288 22
proc/289 22
proc/292 22
proc/294 22
proc/295 22
proc/296 22
proc/298 22
proc/30 22
proc/3100 22
proc/3133 22
proc/31640 22
proc/31644 22
proc/32 22
proc/320 22
proc/321 22
proc/32425 22
proc/32482 22
proc/33 22
proc/37 22
proc/38 22
proc/393 22
proc/394 22
proc/395 22
proc/396 22
proc/402 22
proc/404 22
proc/405 22
proc/406 22
proc/41 22
proc/414 22
proc/42 22
proc/43 22
proc/430 22
proc/431 22
proc/44 22
proc/45 22
proc/47 22
proc/49 22
proc/50 22
proc/51 22
proc/52 22
proc/53 22
proc/54 22
proc/546 22
proc/56 22
proc/57 22
proc/5717 22
proc/5770 22
proc/58 22
proc/5864 22
proc/5877 22
proc/5883 22
proc/59 22
proc/591 22
proc/592 22
proc/593 22
proc/598 22
proc/5982 22
proc/5984 22
proc/599 22
proc/6 22
proc/60 22
proc/600 22
proc/6036 22
proc/6089 22
proc/65 22
proc/656 22
proc/658 22
proc/659 22
proc/66 22
proc/661 22
proc/663 22
proc/67 22
proc/68 22
proc/70 22
proc/7017 22
proc/71 22
proc/710 22
proc/714 22
proc/72 22
proc/73 22
proc/7314 22
proc/74 22
proc/747 22
proc/75 22
proc/753 22
proc/77 22
proc/779 22
proc/78 22
proc/80 22
proc/81 22
proc/83 22
proc/85 22
proc/868 22
proc/88 22
proc/9 22
proc/90 22
proc/91 22
proc/92 22
proc/93 22
proc/95 22
proc/950 22
proc/953 22
proc/96 22
proc/97 22
proc/98 22
proc/980 22
proc/9967 22
proc/9969 22
proc/asound 22
proc/buddyinfo 22
proc/bus 22
proc/cpuinfo 22
proc/devices 22
proc/diskstats 22
proc/dma 22
proc/driver 22
proc/execdomains 22
proc/fb 22
proc/filesystems 22
proc/fs 22
proc/interrupts 22
proc/iomem 22
proc/ipmi 22
proc/kallsyms 22
proc/kcore 22
proc/keys 22
proc/kpageflags 22
proc/latency_stats 22
proc/loadavg 22
proc/locks 22
proc/mdstat 22
proc/misc 22
proc/modules 22
proc/mounts 22
proc/mtrr 22
proc/net 22
proc/pagetypeinfo 22
proc/partitions 22
proc/self 22
proc/slabinfo 22
proc/sys 22
proc/sysrq-trigger 22
proc/sysvipc 22
proc/timer_list 22
proc/uptime 22
proc/version 22
proc/version_signature 22
proc/vmallocinfo 22
proc/vmstat 22
proc/zoneinfo 22
root/* 23
run/acpid.socket 24
run/agentx 24
run/alsa 24
run/avahi-daemon 24
run/bumblebeed.pid 24
run/bumblebee.socket 24
run/crond.reboot 24
run/cups 24
run/dbus 24
run/gunicorn 24
run/initramfs 24
run/kerneloops.pid 24
run/lock 24
run/mlocate.daily.lock 24
run/motd.dynamic 24
run/mount 24
run/network 24
run/network-interface-security 24
run/NetworkManager 24
run/plymouth 24
run/pm-utils 24
run/pppconfig 24
run/rpcbind 24
run/rpcbind.lock 24
run/rpc_pipefs 24
run/rpc.statd.pid 24
run/rsyslogd.pid 24
run/screen 24
run/shm 24
run/sm-notify.pid 24
run/snmpd.pid 24
run/sshd 24
run/sshd.pid 24
run/systemd 24
run/teamviewerd.ipc 24
run/teamviewerd.pid 24
run/udev 24
run/upstart-udev-bridge.pid 24
run/vsftpd 24
run/xrdp 24
sbin/acpi_available 25
sbin/alsa 25
sbin/apparmor_parser 25
sbin/badblocks 25
sbin/blockdev 25
sbin/brltty-setup 25
sbin/capsh 25
sbin/cfdisk 25
sbin/cgdisk 25
sbin/ctrlaltdel 25
sbin/depmod 25
sbin/dhclient 25
sbin/dmsetup 25
sbin/dosfslabel 25
sbin/dumpe2fs 25
sbin/e2image 25
sbin/e2label 25
sbin/e2undo 25
sbin/fatlabel 25
sbin/fdisk 25
sbin/fixparts 25
sbin/fsck 25
sbin/fsck.cramfs 25
sbin/fsck.ext2 25
sbin/fsck.ext3 25
sbin/fsck.ext4 25
sbin/fsck.ext4dev 25
sbin/fsck.fat 25
sbin/fsck.minix 25
sbin/fsck.nfs 25
sbin/fsck.vfat 25
sbin/fsfreeze 25
sbin/fstab-decode 25
sbin/fstrim-all 25
sbin/gdisk 25
sbin/getpcaps 25
sbin/getty 25
sbin/halt 25
sbin/hdparm 25
sbin/hwclock 25
sbin/ifconfig 25
sbin/ifdown 25
sbin/ifquery 25
sbin/init 25
sbin/initctl 25
sbin/insmod 25
sbin/installkernel 25
sbin/ip 25
sbin/ip6tables 25
sbin/ip6tables-apply 25
sbin/ip6tables-restore 25
sbin/ip6tables-save 25
sbin/iptables-restore 25
sbin/iptables-save 25
sbin/iptunnel 25
sbin/isosize 25
sbin/iw 25
sbin/iwconfig 25
sbin/iwevent 25
sbin/iwgetid 25
sbin/iwlist 25
sbin/iwspy 25
sbin/kbdrate 25
sbin/key.dns_resolver 25
sbin/killall5 25
sbin/ldconfig.real 25
sbin/logsave 25
sbin/lsmod 25
sbin/lspcmcia 25
sbin/MAKEDEV 25
sbin/mii-tool 25
sbin/mkfs.bfs 25
sbin/mkfs.cramfs 25
sbin/mkfs.ext3 25
sbin/mkfs.ext4 25
sbin/mkfs.ext4dev 25
sbin/mkfs.fat 25
sbin/mkfs.minix 25
sbin/mkfs.ntfs 25
sbin/mkfs.vfat 25
sbin/mkfs.xfs 25
sbin/mkntfs 25
sbin/mkswap 25
sbin/modinfo 25
sbin/modprobe 25
sbin/mountall 25
sbin/mount.lowntfs-3g 25
sbin/mount.nfs 25
sbin/mount.nfs4 25
sbin/mount.ntfs 25
sbin/nameif 25
sbin/ntfsclone 25
sbin/ntfslabel 25
sbin/ntfsundelete 25
sbin/osd_login 25
sbin/pam_tally2 25
sbin/parted 25
sbin/partprobe 25
sbin/pccardctl 25
sbin/pivot_root 25
sbin/rarp 25
sbin/raw 25
sbin/regdbdump 25
sbin/reload 25
sbin/request-key 25
sbin/resolvconf 25
sbin/restart 25
sbin/rmmod 25
sbin/route 25
sbin/rpcbind 25
sbin/rpc.statd 25
sbin/rtacct 25
sbin/rtmon 25
sbin/runlevel 25
sbin/setvtrgb 25
sbin/sfdisk 25
sbin/sgdisk 25
sbin/shadowconfig 25
sbin/showmount 25
sbin/shutdown 25
sbin/sm-notify 25
sbin/start 25
sbin/startpar 25
sbin/startpar-upstart-inject 25
sbin/status 25
sbin/stop 25
sbin/sulogin 25
sbin/swaplabel 25
sbin/swapon 25
sbin/sysctl 25
sbin/telinit 25
sbin/tune2fs 25
sbin/u-d-c-print-pci-ids 25
sbin/udevadm 25
sbin/udevd 25
sbin/umount.nfs 25
sbin/umount.nfs4 25
sbin/umount.udisks2 25
sbin/upstart-dbus-bridge 25
sbin/upstart-event-bridge 25
sbin/upstart-file-bridge 25
sbin/upstart-local-bridge 25
sbin/upstart-udev-bridge 25
sbin/ureadahead 25
sbin/wpa_action 25
sbin/wpa_supplicant 25
sbin/xtables-multi 25
srv/ftp 26
sys/block 27
sys/bus 27
sys/class 27
sys/firmware 27
sys/fs 27
sys/hypervisor 27
sys/kernel 27
sys/module 27
sys/power 27
usr/bin 29
usr/games 29
usr/include 29
usr/lib 29
usr/lib32 29
usr/NX 29
usr/sbin 29
usr/src 29
var/cache 30
var/crash 30
var/cuda-repo-8-0-local 30
var/games 30
var/lib 30
var/log 30
var/metrics 30
var/opt 30
var/run 30
var/spool 30
var/tmp 30
vmlinuz/* 31
vmlinuz.old/* 32
